# aws/ami/querypie-first-boot.service
[Unit]
Description=Start QueryPie containers on first boot
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service
ConditionPathExists=!/var/lib/querypie/first-boot-done

[Service]
Type=oneshot
WorkingDirectory=/home/ec2-user
ExecStartPre=/bin/mkdir -p /var/lib/querypie
ExecStart=/usr/local/bin/setup.v2.sh --resume
ExecStartPost=/bin/touch /var/lib/querypie/first-boot-done
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
